# deploy/systemd/homelab-monitor-crontab-snapshot.service
[Unit]
Description=Snapshot host user crontabs for homelab-monitor discovery
Documentation=https://github.com/jakekausler/homelab-monitor/blob/main/docs/architecture/cron-logscrape.md
After=cron.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/hm-crontab-snapshot

# --- Hardening ------------------------------------------------------------
# Runs as root (default User) so `crontab -l -u <user>` can read every user's
# crontab via the cron-sanctioned path. Locked down: the ONLY thing this
# service writes is the snapshot directory.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
# ProtectSystem=strict makes the WHOLE filesystem read-only. The snapshot
# script's ONLY write target is /var/lib/homelab-monitor/crontab-snapshot;
# naming the parent /var/lib/homelab-monitor lets the script's `mkdir -p`
# create the snapshot subdir on first run. The script READS
# /var/spool/cron/crontabs and runs /usr/bin/crontab — reads + exec are not
# blocked by ProtectSystem, so they need no entry here.
ReadWritePaths=/var/lib/homelab-monitor
